# deploy/meesho-dashboard-backend.service
[Unit]
Description=Meesho Seller Dashboard — FastAPI backend
After=network.target mongod.service
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/meesho-dashboard/backend
EnvironmentFile=/home/ubuntu/meesho-dashboard/backend/.env
ExecStart=/home/ubuntu/meesho-dashboard/backend/.venv/bin/uvicorn server:app --host 127.0.0.1 --port 8001 --workers 2
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
